# systemd/ctrl-hold-voice.service
[Unit]
Description=Voice Input Daemon (Ctrl hold + Whisper GPU)
After=graphical-session.target

[Service]
ExecStart=__PYTHON__ __DAEMON__
Restart=always
RestartSec=5
Environment=DISPLAY=:1
Environment=XAUTHORITY=__XAUTHORITY__
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/__UID__/bus
Environment=VOICE_MIC_TARGET=__MIC_TARGET__
Environment=VOICE_LLM_URL=http://127.0.0.1:8000/v1/chat/completions
Environment=VOICE_LLM_ENABLED=1
Environment=VOICE_HOLD_SECONDS=0.5
Environment=VOICE_WHISPER_MODEL=medium
Environment=VOICE_WHISPER_DEVICE=cuda
Environment=VOICE_BEAM_SIZE=1
Environment=VOICE_WHISPER_LANGUAGE=zh

[Install]
WantedBy=default.target
